# systemd/user/battery-notifications.service
[Unit]
Description=Battery notifications

[Service]
Type=oneshot
ExecStart=/usr/bin/battery-notifications.sh

[Install]
WantedBy=graphical.target

